# web/serial_scanner/memcached.conf
[program:memcached]
command=/usr/bin/memcached
autostart=true
autorestart=true
user=www-data